# debian/unfold.postinst
#!/bin/bash
# if this requires a service to be running, add something like this
# update-rc.d unfold defaults
