# debian/unfold.postinst
#!/bin/bash
# if this requires a service to be running, add something like this
# update-rc.d unfold defaults
/usr/share/unfold/manage.py syncdb
/usr/share/unfold/manage.py migrate
chmod 777 /usr/share/unfold/myslice.sqlite3
chmod 777 /usr/share/unfold
a2dissite default
a2ensite myslice.conf
service apache2 restart
